# etc/systemd/system/crm-repeats.service
# /etc/systemd/system/crm-repeats.service
# Запускает проверку повторных сделок через API CRM

[Unit]
Description=GrassCRM — проверка повторных сделок
After=network.target crm.service
Requires=crm.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/curl -sf \
    -H "X-Internal-API-Key: ${INTERNAL_API_KEY}" \
    http://127.0.0.1:8000/api/service/run-repeats
EnvironmentFile=/etc/crm/.env
StandardOutput=journal
StandardError=journal
